# .config/systemd/user/mosh.service
[Unit]
Description=Mosh Service

[Service]
ExecStart=/usr/bin/mosh-server

[Install]
WantedBy=multi-user.target
